# infra/ansible/templates/aslf-super.service.j2
[Unit]
Description=ASLF-OSINT Super Control API
After=network.target

[Service]
Type=simple
WorkingDirectory={{ aslf_install_dir }}
EnvironmentFile={{ aslf_install_dir }}/.env.super
ExecStart={{ aslf_install_dir }}/.venv/bin/python {{ aslf_install_dir }}/run_control_api.py --host 0.0.0.0 --port {{ aslf_api_port }}
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

